# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=/bin/bash -c "java -jar ./target/mock-server-jar-with-dependencies.jar"
username=root
autorestart=false
redirect_stderr=true